# deploy/supervisor.conf
; ============================================================
; crmoffice — Supervisor Configuration
; Target: Ubuntu 22.04 / 24.04
; ============================================================
; Deploy:
;   cp deploy/supervisor.conf /etc/supervisor/conf.d/crmoffice.conf
;   supervisorctl reread
;   supervisorctl update
;   supervisorctl start crmoffice-scheduler:*
;   supervisorctl start crmoffice-worker-default:*
;   supervisorctl start crmoffice-worker-high:*
;   supervisorctl start crmoffice-horizon:*
; ============================================================

; ---- Group: crmoffice (manage all together) ----
[group:crmoffice]
programs=crmoffice-scheduler,crmoffice-worker-default,crmoffice-worker-high,crmoffice-horizon

; ============================================================
; 1. Laravel Scheduler (runs artisan schedule:run every minute)
; ============================================================
[program:crmoffice-scheduler]
process_name=%(program_name)s
command=php /var/www/crmoffice/artisan schedule:work
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/crmoffice-scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/crmoffice-scheduler-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stopwaitsecs=600

; ============================================================
; 2. Default Queue Workers (3 processes)
; ============================================================
[program:crmoffice-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/crmoffice/artisan queue:work --queue=default --sleep=3 --tries=3 --max-time=3600 --timeout=60
autostart=true
autorestart=true
user=www-data
numprocs=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/crmoffice-worker-default.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/crmoffice-worker-default-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stopwaitsecs=3600

; ============================================================
; 3. High Priority Queue Worker (1 process)
; ============================================================
[program:crmoffice-worker-high]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/crmoffice/artisan queue:work --queue=high --sleep=1 --tries=5 --max-time=3600 --timeout=120
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/crmoffice-worker-high.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/crmoffice-worker-high-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stopwaitsecs=3600

; ============================================================
; 4. Laravel Horizon (Redis-backed queue monitoring dashboard)
; ============================================================
[program:crmoffice-horizon]
process_name=%(program_name)s
command=php /var/www/crmoffice/artisan horizon
autostart=true
autorestart=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/crmoffice-horizon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/crmoffice-horizon-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
stopwaitsecs=3600

; ============================================================
; Log rotation (daily rotation via logrotate.d)
; ============================================================
; Also create /etc/logrotate.d/crmoffice-supervisor:
;
;   /var/log/supervisor/crmoffice-*.log {
;       daily
;       missingok
;       rotate 14
;       compress
;       delaycompress
;       notifempty
;       copytruncate
;       su www-data www-data
;   }
